# gdm.conf
description "Gnome Display Manager"

start on (runlevel [5] and started dbus and started hal)
stop on runlevel [!5]

respawn
expect fork

exec /usr/sbin/gdm
