# rc-shutdown.conf
start on runlevel [06]
stop on runlevel [!06]

task

console owner
exec /etc/rc.shutdown
