# udevmonitor.conf
description "log initial device creation"

start on (startup and starting udevtrigger)
stop on stopped udevtrigger

exec /sbin/udevadm monitor -e >/dev/.udev.log
