# udevtrigger.conf
description "cold plug devices"

start on started udev

task

exec /sbin/udevadm trigger
