# udevtrigger.conf
description "cold plug devices"

start on (startup and started udev)

task

exec /sbin/udevadm trigger
